# configs/baseline/airootfs/etc/systemd/system/ollama.service.d/override.conf
[Service]
Environment="HOME=/var/lib/ollama" "GIN_MODE=release" "OLLAMA_MODELS=/usr/share/ollama/.ollama/models"
